# install/temp-listener.service
[Unit]
Description=temp-report email reply listener
After=network.target
StartLimitIntervalSec=10

[Service]
Type=simple
Restart=always
RestartSec=1
User=pi
ExecStart=/usr/bin/env python3 /home/pi/temp-report/listener.py

[Install]
WantedBy=multi-user.target
